# eth-supervisord-master.conf
[supervisord]
logfile=eth-supervisord.log

[program:eth]
stdout_logfile=eth.log
stderr_logfile=eth.err
command=eth -c "Develop Server Eight" -m 0