# eth-supervisord.conf
[supervisord]
logfile=eth-supervisord.log

[program:eth]
stdout_logfile=eth.log
stderr_logfile=eth.err
command=eth -c buildslave -m 0 -r poc-9.ethdev.com