# eth-go-supervisord.conf
[supervisord]
logfile=eth-go-supervisord.log

[program:eth]
stdout_logfile=ethereum.log
stderr_logfile=ethereum.err
command=ethereum -id=buildslave