# startup/eth-supervisord-master.conf
[supervisord]
logfile=eth-supervisord.log

[program:eth]
stdout_logfile=eth.log
stderr_logfile=eth.err
command=eth -c "Develop Server Eight" -m off --upnp off --listen-ip 5.1.83.225
