# startup/eth-supervisord.conf
[supervisord]
logfile=eth-supervisord.log

[program:eth]
stdout_logfile=eth.log
stderr_logfile=eth.err
command=eth --client-name buildslave --mining off --bootstrap --master 0
