# startup/geth-supervisord.conf
[supervisord]
logfile=geth-supervisord.log

[program:eth]
stdout_logfile=geth.log
stderr_logfile=geth.err
command=geth
